# jsonnet_examples/ini_files/supervisord_conf.ini
[program:flask_app_api]
autorestart = true
autostart = true
command = /usr/bin/python /home/fedora/apps/flask_app_api.py
directory = /home/fedora/apps/
stderr_logfile = /home/fedora/apps/logs/flas_app_api.err.log
stdout_logfile = /home/fedora/apps/logs/flas_app_api.out.log
[program:flask_app_backend]
autorestart = true
autostart = true
command = /usr/bin/python /home/fedora/apps/flask_app_backend.py
directory = /home/fedora/apps/
stderr_logfile = /home/fedora/apps/logs/flas_app_backend.err.log
stdout_logfile = /home/fedora/apps/logs/flas_app_backend.out.log
[program:flask_app_frontend]
autorestart = true
autostart = true
command = /usr/bin/python /home/fedora/apps/flask_app_frontend.py
directory = /home/fedora/apps/
stderr_logfile = /home/fedora/apps/logs/flas_app_frontend.err.log
stdout_logfile = /home/fedora/apps/logs/flas_app_frontend.out.log
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
[supervisord]
directory = /tmp
logfile = /tmp/supervisord.log
logfile_backups = 5
logfile_maxbytes = 75MB
loglevel = info
pidfile = /tmp/supervisord.pid
user = fedora
